# display/schedule.crontab
0 0 * * * /home/webserver/bin/display #clear
0 5 * * * /home/webserver/bin/display /home/webserver/.local/share/eink/5.bin
0 12 * * * /home/webserver/bin/display /home/webserver/.local/share/eink/12.bin
0 18 * * * /home/webserver/bin/display /home/webserver/.local/share/eink/18.bin
